# templates/cloudflared.service.j2
[Unit]
Description=cloudflared tunnel for %I
After=network.target

[Service]
Type=simple
ExecStart=cloudflared tunnel --config {{ cf_config_dir_cf_named_tunnels }}/%i.yml run %i
User={{ cf_systemd_user }}
Group={{ cf_systemd_group }}

Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target